# API/systemd/timers/gold_price_daily_percentage.timer
[Unit]
Description=Runs Gold Price Daily Percentage Daily

[Timer]
OnCalendar=*-*-* 07:05:00
Persistent=true

[Install]
WantedBy=timers.target
